# ansible/pi-server/roles/rclone_backup/templates/docker-backup.timer.j2
[Unit]
Description=timer to call backupscript

[Timer]
OnCalendar={{ rclone_backup_timer_oncalendar }}

[Install]
